# services/rpcbind/run
#!/command/execlineb -P

exec rpcbind -f